# src/dist/pad.service
[Unit]
Description=Starts pad

[Service]
Type=simple
Restart=always
ExecStart=/opt/pad/bin/pad
KillMode=process
TimeoutSec=180

[Install]
WantedBy=network.target